# x-bot-blocker.service
[Unit]
Description=X Bot Blocker Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/x-bot-blocker
Environment=PYTHONPATH=/home/ubuntu/x-bot-blocker
ExecStart=/home/ubuntu/x-bot-blocker/start_bot.sh
ExecStop=/home/ubuntu/x-bot-blocker/stop_bot.sh
TimeoutStopSec=30
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 